# ops/systemd/hwistock-api.service
# hwiStock API service template (UNIT-002). Do not execute from repo checkout.
# Install/copy to /etc/systemd/system/ only after operator review.

[Unit]
Description=hwiStock FastAPI health and read-only runner status API
After=network.target
Wants=postgresql.service

[Service]
Type=simple
User=hwi
Group=hwi
WorkingDirectory=/data/workspace/My/hwiStock/backend
EnvironmentFile=/home/hwi/.config/hwistock/hwistock.env
Environment=HWISTOCK_BIND_HOST=127.0.0.1
ExecStart=/usr/bin/env python3 /data/workspace/My/hwiStock/ops/systemd/run_with_daily_logs.py --name hwistock-api -- /usr/bin/python3 -m uvicorn server:app --host 127.0.0.1 --port 5001
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
