# ops/systemd/user/hwistock-intel-collector.service
# hwiStock 24h market intelligence collector tick.
# Metadata-only news/disclosure ingestion; no broker/order routing.

[Unit]
Description=hwiStock market intelligence collector tick
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/data/workspace/My/hwiStock
EnvironmentFile=-/home/hwi/.config/hwistock/hwistock.env
EnvironmentFile=-/home/hwi/.config/hwistock/runtime-mode.env
ExecStart=/usr/bin/env python3 ops/systemd/run_with_daily_logs.py --name hwistock-intel-collector -- /usr/bin/env bash -lc 'source ./env.sh >/dev/null 2>&1 && source ops/systemd/load_runtime_mode_env.sh && python backend/service/market_intelligence_ingestion.py --once'
NoNewPrivileges=true
