# ops/systemd/user/hwistock-api.service
# hwiStock local API user service.
# Local-only FastAPI backend on 127.0.0.1:5001.

[Unit]
Description=hwiStock local FastAPI API (user, loopback)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/data/workspace/My/hwiStock
ExecStartPre=/usr/bin/env mkdir -p /data/workspace/My/hwiStock/logs/systemd
StandardOutput=append:/data/workspace/My/hwiStock/logs/systemd/hwistock-api.log
StandardError=append:/data/workspace/My/hwiStock/logs/systemd/hwistock-api.err.log
EnvironmentFile=-/home/hwi/.config/hwistock/hwistock.env
Environment=HWISTOCK_BIND_HOST=127.0.0.1
ExecStart=/usr/bin/env bash -lc 'source ./env.sh >/dev/null 2>&1 && cd backend && exec uvicorn server:app --host 127.0.0.1 --port 5001'
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
