# ops/systemd/user/hwistock-ai-analysis.timer
# hwiStock DeepSeek/MoonBridge analysis timer.

[Unit]
Description=hwiStock DeepSeek analysis every hour

[Timer]
OnCalendar=*-*-* *:00:00
AccuracySec=1min
Persistent=true
Unit=hwistock-ai-analysis.service

[Install]
WantedBy=timers.target
