# ops/systemd/user/hwistock-api.service
# hwiStock local API user service.
# Local-only FastAPI backend on 127.0.0.1:5001.

[Unit]
Description=hwiStock local FastAPI API (user, loopback)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/data/workspace/My/hwiStock
EnvironmentFile=-/home/hwi/.config/hwistock/hwistock.env
EnvironmentFile=-/home/hwi/.config/hwistock/runtime-mode.env
Environment=HWISTOCK_BIND_HOST=127.0.0.1
Environment=HWISTOCK_MARKET_DATA_SOURCE=kis_market_mode_aware
Environment=HWISTOCK_CALENDAR_PATH=/data/workspace/My/hwiStock/config/market-calendar/krx-nxt-trading-days.json
Environment=HWISTOCK_RUNTIME_CALENDAR_PATH=/home/hwi/.config/hwistock/market-calendar/krx-nxt-trading-days.json
Environment=HWISTOCK_ALLOW_WEEKDAY_CALENDAR_FALLBACK=false
ExecStart=/usr/bin/env python3 ops/systemd/run_with_daily_logs.py --name hwistock-api -- /usr/bin/env bash -lc 'source ./env.sh >/dev/null 2>&1 && source ops/systemd/load_runtime_mode_env.sh && cd backend && exec uvicorn server:app --host 127.0.0.1 --port 5001'
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
