# ops/systemd/user/hwistock-kis-paper-health.service
# KIS paper/mock read-health tick.
# Paper/mock only. No order endpoint.

[Unit]
Description=hwiStock KIS paper/mock read-health tick
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/data/workspace/My/hwiStock
EnvironmentFile=-/home/hwi/.config/hwistock/hwistockApi.env
ExecStart=/usr/bin/env bash -lc 'source ./env.sh >/dev/null 2>&1 && python backend/service/kis_paper_health.py --once'
NoNewPrivileges=true
