# ops/systemd/user/hwistock-kis-market-data.timer
# hwiStock KIS market-data collector timer.

[Unit]
Description=hwiStock KIS mode-aware market-data collector every minute

[Timer]
OnBootSec=90s
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=hwistock-kis-market-data.service

[Install]
WantedBy=timers.target
