# ops/systemd/user/hwistock-intel-collector.timer
# hwiStock 24h market intelligence collector timer.

[Unit]
Description=hwiStock market intelligence collector every 10 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=hwistock-intel-collector.service

[Install]
WantedBy=timers.target
