# ops/systemd/user/hwistock-ai-flash.timer
# hwiStock DeepSeek Flash 10-minute trade-document timer.

[Unit]
Description=hwiStock DeepSeek Flash trade document during paper KRX decision window

[Timer]
OnCalendar=Mon..Fri *-*-* 09..14:00/10:00
AccuracySec=30s
Persistent=false
Unit=hwistock-ai-flash.service

[Install]
WantedBy=timers.target
