# ops/systemd/user/hwistock-frontend.service
# hwiStock local dashboard frontend user service.
# Local-only Next.js dashboard on 127.0.0.1:5000.

[Unit]
Description=hwiStock local dashboard frontend (user, loopback)
After=network-online.target hwistock-api.service
Wants=network-online.target hwistock-api.service

[Service]
Type=simple
WorkingDirectory=/data/workspace/My/hwiStock
Environment=HOSTNAME=127.0.0.1
Environment=PORT=5000
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
ExecStart=/usr/bin/env python3 ops/systemd/run_with_daily_logs.py --name hwistock-frontend -- /usr/bin/env bash -lc 'source ./env.sh >/dev/null 2>&1 && cd frontend-web && exec pnpm start'
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
