# ops/systemd/user/hwistock-ai-flash.service
# hwiStock DeepSeek Flash 10-minute trade-document tick.
# Report-only; AI artifacts never submit broker orders.

[Unit]
Description=hwiStock DeepSeek Flash trade-document tick
After=network-online.target hwistock-ai-analysis.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/data/workspace/My/hwiStock
EnvironmentFile=-/home/hwi/.config/hwistock/deepseek.env
EnvironmentFile=-/home/hwi/.config/hwistock/runtime-mode.env
Environment=HWISTOCK_DEEPSEEK_MODEL=deepseek-v4-flash
Environment=HWISTOCK_AI_EVENT_LIMIT=20
Environment=HWISTOCK_AI_REASONING_EFFORT=low
Environment=HWISTOCK_RUNTIME_CALENDAR_PATH=/home/hwi/.config/hwistock/market-calendar/krx-nxt-trading-days.json
ExecStart=/usr/bin/env python3 ops/systemd/run_with_daily_logs.py --name hwistock-ai-flash -- /usr/bin/env bash -lc 'source ./env.sh >/dev/null 2>&1 && source ops/systemd/load_runtime_mode_env.sh && python ops/systemd/ensure_market_calendar.py --investment-mode "${HWISTOCK_INVESTMENT_MODE:-paper}" --date today --horizon-days 14 --write --write-evidence && python backend/service/ai_analysis_runner.py --once --job flash-10m --model "${HWISTOCK_DEEPSEEK_MODEL:-deepseek-v4-flash}"'
NoNewPrivileges=true
